# kiosk_clock.service
[Unit]
Description=Kiosk Clock - Digital Clock with Calendar Integration
After=graphical-session.target
Wants=graphical-session.target
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/kiosk_clock
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStartPre=/bin/sleep 10
ExecStart=/home/pi/kiosk_clock/start_kiosk_clock.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/kiosk_clock

[Install]
WantedBy=graphical-session.target 